# proxy.conf
[program:proxy]
command=bash -c /usr/src/app/proxy.sh
user=root